# nurds.container
[Unit]
Description=The nurds container!
After=network.target

[Container]
Image=ghcr.io/geraveoyomama/nurds-container:latest
ContainerName=nurds
AutoUpdate=registry

PublishPort=8080:80

LogDriver=journald

[Service]
# Inform systemd of additional exit status
SuccessExitStatus=0 143
Restart=always

[Install]
# Start by default on boot
WantedBy=multi-user.target default.target
